# server/systemd/device-policy-manager.service.in
[Unit]
Description=Device Policy Manager
Wants=lbs-server.service net-config.service pulseaudio.service
After=lbs-server.service net-config.service pulseaudio.service

[Service]
Type=simple
User=security_fw
Group=security_fw
SmackProcessLabel=System
ExecStartPre=${CMAKE_INSTALL_PREFIX}/bin/dpm-storage-builder
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/device-policy-manager
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE

[Install]
WantedBy=multi-user.target
