# systemd/auth-fw.service.in
[Unit]
Description=Start @SERVICE_NAME@ server
Requires=@SOCK_PASSWD_CHECK@ @SOCK_PASSWD_SET@ @SOCK_PASSWD_RESET@ @SOCK_PASSWD_POLICY@

[Service]
Type=notify
ExecStart=@BIN_DIR@/@SERVICE_NAME@
User=security_fw
Group=security_fw
SmackProcessLabel=System

[Install]
WantedBy=multi-user.target
