# systemd/cert-server.service.in
[Unit]
Description=cert server for cert-svc managing root ssl certs
Wants=cert-server.socket

[Service]
ExecStart=@TZ_SYS_BIN@/cert-server
User=@USER_NAME@
Group=@GROUP_NAME@
SmackProcessLabel=@SMACK_DOMAIN_NAME@
Sockets=cert-server.socket
