# cert-svc/systemd/cert-server.service.in
[Unit]
Description=cert server for cert-svc managing root ssl certs

[Service]
ExecStart=@BIN_DIR@/cert-server
User=@USER_NAME@
Group=@GROUP_NAME@
SmackProcessLabel=@SMACK_DOMAIN_NAME@
