# cert-svc/systemd/cert-server.socket.in
[Socket]
ListenStream=@SERVER_STREAM@
SocketUser=@USER_NAME@
SocketGroup=@GROUP_NAME@
SocketMode=0777

[Install]
WantedBy=sockets.target
