# device-policy-manager/server/systemd/device-policy-manager.service.in
[Unit]
Description=Device Policy Manager
Wants=pulseaudio.service
After=pulseaudio.service

[Service]
User=${USER_NAME}
Group=${GROUP_NAME}
SmackProcessLabel=${SMACK_LABEL}
ExecStartPre=${CMAKE_INSTALL_PREFIX}/bin/dpm-storage-builder
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/device-policy-manager -o 5000
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE

[Install]
WantedBy=multi-user.target
