# auth-fw/systemd/auth-fw-passwd-check.socket.in
[Socket]
ListenStream=@RUN_DIR@/.@SOCK_PASSWD_CHECK@
SocketUser=security_fw
SocketGroup=security_fw
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@
Service=@SERVICE_NAME@.service

[Install]
WantedBy=sockets.target
