# cert-checker/systemd/cert-checker-popup.service.in
[Unit]
Description=Cert-checker popup service

[Service]
EnvironmentFile=@POPUP_ENV_PATH@
ExecStart=@BIN_DIR@/@SERVICE_NAME@-popup
