# cert-checker/systemd/cert-checker-popup.socket.in
[Unit]
Description=Cert-checker popup socket

[Socket]
ListenStream=@POPUP_STREAM@
SocketMode=0777

[Install]
WantedBy=sockets.target
